# raspi/nfc/nfc-api.service
[Unit]
Description=NFC API Server for LiVraria
After=network.target pcscd.service
Requires=pcscd.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/livraria/LiVraria/raspi/nfc
ExecStart=/home/livraria/LiVraria/raspi/nfc/.venv/bin/python nfc_api_server.py
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
